# roles/base/files/CentOS/7/etc/profile.d/gcb.sh
# TWGCB-01-008-0238
readonly TMOUT=900 ; export TMOUT
umask 027